# nkelemen18/linux/roles/restic/templates/systemd-restic-backup.timer.j2
[Unit]
Description=Restic backup ({{ restic_schedule }})

[Timer]
OnCalendar={{ restic_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
