# nkelemen18/linux/roles/restic/templates/systemd-restic-cleanup.timer.j2
[Unit]
Description=Cleanup restic repository ({{ restic_cleanup_schedule }})

[Timer]
OnCalendar={{ restic_cleanup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
